# wService.conf
[program:WeatherService]
directory=/root/Weather-Service/
command=python service.py
autostart=true
autorestart=true

